# systemd/system/pihole_curfew_start.service
[Unit]
Description=pihole curfew start
[Service]
ExecStart=pihole-block curfew
