# systemd/system/pihole_curfew_stop.service
[Unit]
Description=pihole curfew stop
[Service]
ExecStart=pihole-unblock curfew
